# authorization-api/s6-rc.d/guardian_authorization_api/run
#!/command/with-contenv bash

cd $SERVICE_DIR
exec s6-setuidgid $SERVICE_USER gunicorn --worker-class uvicorn.workers.UvicornWorker --bind 0.0.0.0:8000 guardian_authorization_api.main:app
