# authorization-api/s6-rc.d/guardian_authorization_api/run
#!/bin/sh

exec s6-setuidgid guardian with-contenv gunicorn --worker-class uvicorn.workers.UvicornWorker --bind 0.0.0.0:8000 guardian_authorization_api.main:app