# opa/s6-rc.d/opa/run
#!/command/with-contenv bash

cd $SERVICE_DIR
exec s6-setuidgid $SERVICE_USER opa run --disable-telemetry --server -c opa_config.yaml --addr 0.0.0.0:8181 -b BaseBundle
